# examples/marieluise/chirpsounder_kho_archive_plots.service
[Unit]
Description=KHO archived ionosonde web quick-look plots
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/src/chirpsounder2
ExecStart=%h/src/chirpsounder2/.venv/bin/python %h/src/chirpsounder2/plot_archive_quicklooks.py --config examples/marieluise/kho_archive.ini --data-dir /mnt/shovel/ionosonde --web-dir /var/www/html/iono
Restart=always
RestartSec=30
Nice=10

[Install]
WantedBy=default.target
